# openbox/.config/autostart/panel.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=tint2
Name[en_US]=tint2
Icon=krfb
Exec=tint2
Comment[en_US]=
StartupNotify=true
